# ushine/files/etc/uci-defaults/60_ushine_rename_root_to_admin
uci set rpcd.@login[0].username='admin'
uci commit rpcd
