# ushine/files/etc/uci-defaults/61_ushine_system_timezone
uci set system.@system[-1].timezone="CST-8"
uci set system.@system[-1].zonename=""Asia/Taipei
uci commit system
